# files/etc/cron.hourly/journal
#!/bin/bash

cd /home/kirill/sync/data/appdata/myJournal
/usr/bin/git add .
/usr/bin/git commit -m . >> ~/1.txt
