# files/etc/cron.hourly/journal
#!/bin/bash

cd /home/kirill/sync/data/appdata/myJournal
git add .
git commit -m .
